# config/supervisor/foodcenter.conf
[program:foodcenter]

command     = spawn-fcgi -a 127.0.0.1 -p 9002 -d /var/www/foodcenter -f /var/www/foodcenter/index.wsgi
directory   = /var/www/foodcenter

